# image/etc/periodic/15min/calibre-import
#!/bin/sh

/opt/calibre/calibredb add /calibre-import/ -r --with-library http://localhost:8080/#calibre-library
rm /calibre-import/*